# stow/home/etc/sv/kmonad/run
#!/bin/sh

exec 2>&1
exec /usr/bin/kmonad /home/mamapitufo/.config/kmonad/internal.kbd || exit 1
