# stow/yibit/etc/sv/runsvdir-mamapitufo/run
#!/bin/sh

chpst -u mamapitufo mkdir -p -m 0700 /run/user/$(id -u mamapitufo)/runit || exit 1
exec chpst -u "mamapitufo:$(id -Gn mamapitufo | tr ' ' ':')" runsvdir /home/mamapitufo/service
