# stow/yibit/etc/sv/runsvdir-mamapitufo/run
#!/bin/sh

exec 2>&1
chpst -u mamapitufo mkdir -p -m 0700 /run/user/$(id -u mamapitufo)/runit
exec chpst -u "mamapitufo:$(id -Gn mamapitufo | tr ' ' ':')" runsvdir /home/mamapitufo/service
