# stow/yibit/etc/sv/runsvdir-mamapitufo/run
#!/bin/sh

exec chpst -u "mamapitufo:$(id -Gn mamapitufo | tr ' ' ':')" runsvdir /home/mamapitufo/service
